# roles/uplink/files/etc/udev/rules.d/99-filter.rules
ACTION=="add", SUBSYSTEM=="net", \
    RUN+="/usr/bin/bash /usr/local/bin/uplink-filter.sh '$env{ACTION}' '$name'" 
