# roles/bluetooth_sink/files/etc/systemd/system/bluealsa.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa --profile=a2dp-sink --a2dp-force-audio-cd

[Install]
WantedBy=bluetooth.service
