# roles/mpdutils/templates/etc/systemd/system/playlistbot.service.j2
[Unit]
Description=Manipulate mpd playlists

[Service]
Type=oneshot
ExecStart=sudo -u nobody /usr/local/bin/playlistbot.sh {{ mpdutils_playlistbot_options }}
